# uf_scrapy/crontab
* * 9 * * root . /root/envs.sh;cd /app && python /app/run.py >> /var/log/cron.log 2>&1
